# BackEnd/deploy/daphne.conf
[program:daphne]
user = root
command = /app/deploy/start_daphne.sh
directory=/app/core
autostart=true
autorestart=true
stdout_logfile = /logs/daphne.log
redirect_stderr = true